# deploy/supervisor.conf
[program:faucet-frontend]
command=/home/ubuntu/.bun/bin/bun start
directory=/home/ubuntu/faucet/frontend
user=ubuntu
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/faucet/frontend.err.log
stdout_logfile=/var/log/faucet/frontend.out.log
environment=NODE_ENV=production,PATH="/home/ubuntu/.bun/bin:/usr/local/bin:/usr/bin:/bin"
stopasgroup=true
killasgroup=true
redirect_stderr=true
